# templates/LaunchAgents/com.fabs.mortadela.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.fabs.mortadela</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/<%= user %>/bin/mortadela</string>
      <string>/Users/<%= user %>/gdrive</string>
    </array>
    <key>WatchPaths</key>
    <array>
      <string>/Users/<%= user %>/gdrive</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>21600</integer>
  </dict>
</plist>
